# udev/rules.d/00-autorandr.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/su sebastian -c 'DISPLAY=:0 autorandr -c'"
