# systemd/system/hdparm.service
[Unit]
Description=Disable harmful head parking

[Service]
ExecStart=/sbin/hdparm -B 255 /dev/sda

[Install]
WantedBy=basic.target
