# systemd/system/slim.service
.include /usr/lib/systemd/system/slim.service

[Service]
StandardOutput=null
