# cron.daily/btrfs-backup
#! /bin/bash
/usr/local/bin/snapshot daily 7
