# cron.daily/btrfs-backup
#! /bin/bash
/etc/btrfs-snapshot /mnt/real-root daily 7
/etc/btrfs-snapshot /mnt/real-root/home daily 7
