# cron.weekly/btrfs-backup
#! /bin/bash
/etc/btrfs-snaphot /mnt/real-root weekly 100
/etc/btrfs-snaphot /mnt/real-root/home weekly 1000
