# cron.weekly/btrfs-backup
#! /bin/bash
/usr/local/bin/snapshot weekly 100
