# cron.d/atop
# start atop daily at midnight
0 0 * * * /etc/atop/atop.daily
